# disable_bd_prochot.service
[Unit]
Description=Disable BD PROCHOT

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disable_bd_prochot.sh

[Install]
WantedBy=multi-user.target
